# fpc.service
[Unit]
Description=FPC Service

[Service]
ExecStart=/scripts/start_pcap_final.sh

[Install]
WantedBy=multi-user.target
