# fpc.conf
description	"conf to start fpc at startup"
start on startup
task
exec /nsm/scripts/start_pcap_final.sh
