# fpc.conf
description	"conf to start fpc at startup"
start on startup
task
exec /scripts/start_pcap.sh
